# .config/systemd/user/vim_cleanup_undofiles.service
[Unit]
Description=Purge old vim undofiles

[Service]
Type=oneshot
ExecStart=/bin/bash $HOME/.scripts/vim_cleanup_undofiles
